# services/desktop/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv sh
set -e

: "${PORT:=3000}"

# Ensure nginx runtime dir exists and disable default site that may bind to :80
mkdir -p /run/nginx || true
rm -f /etc/nginx/sites-enabled/default /etc/nginx/sites-available/default 2>/dev/null || true

# Render nginx.conf from template with only $PORT substituted
/usr/bin/envsubst '$PORT' < /etc/nginx/nginx.conf.tmpl > /etc/nginx/nginx.conf

exec /usr/sbin/nginx -g 'daemon off;'
